# TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/base-filesystem.list
/etc/openwrt_version
/sbin/sound
/etc/init.d/done
/sbin/network_get_firm
/etc/hotplug2-init.rules
/etc/preinit
/sbin/mount_root
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/etc/shells
/etc/protocols
/sbin/loadprofile
/usr/sbin/dut_auto_reboot
/etc/hotplug.d/button/00-button
/etc/config/system
/sbin/tddp_hook_close_txbf
/etc/sysctl.conf
/sbin/sysupgrade
/lib/upgrade/keep.d/base-filesystem
/sbin/led.sh
/etc/banner
/etc/init.d/auto-reboot
/etc/shadow
/lib/usb/usb.sh
/etc/rc.common
/sbin/gratuitous_arp
/etc/init.d/watchdog
/etc/init.d/sysntpd
/etc/rc.local
/sbin/wifi_check_country
/lib/functions/uci-defaults.sh
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/init.d/usb
/usr/sbin/timeout_reboot_f
/etc/services
/sbin/wifi_commit
/etc/init.d/garp
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/etc/group
/etc/uci-defaults/migrate-sysctl
/etc/functions.sh
/sbin/firstboot
/bin/login.sh
/sbin/loadconfig
/etc/uci-defaults/migrate-shadow
/etc/hosts
/sbin/gen_passwd
/lib/functions.sh
/etc/config/network
/sbin/sysdebug
/etc/init.d/rcS
/lib/usb/usblog.sh
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/etc/profile
/lib/functions/uci-protocol-defaults.sh
/lib/functions/service.sh
/etc/init.d/umount
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/sbin/getfirm
/etc/init.d/boot
/lib/functions/gpio-button-hotplug.sh
/lib/functions/boot.sh
/bin/ipcalc.sh
/etc/diag.sh
