# TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/product-file.list
/etc/init.d/bcmtm
/sbin/board_restore
/lib/preinit/06_init_insmod_usbhost_scsi
/lib/wifi/brcm_wlan_config.txt
/usr/sbin/wps_status
/sbin/board_factory
/etc/hotplug.d/switch/65-iptv
/sbin/network_firm
/usr/lib/.gitignore
/usr/sbin/status
/etc/hotplug.d/iface/65-iptv
/lib/iptv/iptv_config.lua
/lib/switch/core_phy.sh
/etc/hotplug.d/net/10-wl_cpu_affinity
/lib/bcmenet/bcmenet_c5400s.sh
/lib/wifi/brcmwifi.sh
/sbin/switch_link_down
/lib/bcmenet/bcmenet_4908.sh
/etc/init.d/wireless
/etc/init.d/bcmenet
/sbin/wifi_button
/sbin/portset
/lib/print_server/core_print_server.sh
/etc/hotplug.d/net/21-switch_ports
/lib/iptv/iptv_network.sh
/lib/wifi/wps_button
/sbin/network_lanwan_swap
/lib/led_switch
/etc/hotplug.d/firewall/65-iptv
/sbin/wifi_firm
/lib/bcmenet/bcmenet_rtl8367s.sh
/etc/hotplug.d/net/20-client_mgmt
/etc/hotplug.d/usb/21-usb_reset
/lib/preinit/05_init_interfaces_brcm
/lib/switch/core.sh
/lib/preinit/71_init_etc_ramfs
/lib/iptv/iptv_func.sh
/lib/print_server/printer_driver.sh
/usr/sbin/lan_connected
/etc/init.d/bcm_cpu_affinity
/etc/init.d/iptv
/etc/hotplug.d/block/60-enable_gro
/etc/init.d/switch
/etc/hotplug.d/net/01-wds
/lib/nat/nat_hw.sh
/etc/hotplug.d/net/22-add_neigh
/sbin/portspeed
/etc/init.d/phyport
/usr/sbin/wan_connected
/etc/hotplug.d/iface/02-wireless-vlan
/lib/iptv/iptv_core.sh
/sbin/switch_link_up
/etc/hotplug.d/usb/20-usb_cpu_affinity
