# TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/sync-server.list
/lib/sync-server/scripts/sync_wifi
/usr/lib/lua/sync-script.lua
/lib/sync-server/scripts/trans_backup_wcfg
/lib/sync-server/scripts/trans_main_wcfg
/usr/bin/sync-server
/lib/sync-server/scripts/request_clients
/etc/init.d/sync-server
/lib/sync-server/scripts/request
/etc/init.d/generate_group_id
