# TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/time-machine.list
/lib/time_machine/tm_usb_remove.sh
/etc/init.d/time_machine
/etc/hotplug.d/usb/70-time_machine
/lib/time_machine/tm_usb_add.sh
/etc/hotplug.d/iface/70-time_machine
/usr/sbin/tm_getpath
/lib/time_machine/tm_ifup.sh
/usr/sbin/tm_checkrun
