# TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/wportalctrl.list
/usr/bin/wportalctrl
/etc/init.d/wportal
/etc/hotplug.d/iface/01-wportal
/lib/wportal/clear_wan_error.lua
/lib/wportal/saveconfig.lua
/lib/wportal/clear_upgrade.lua
/lib/wportal/wportaltimecmd.sh
/lib/wportal/wportaltime.sh.tmp
/etc/hotplug.d/iface/99-wportal
/lib/wportal/wportal.sh
