# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/etc/hotplug.d/firewall/66-ffs
#!/bin/sh

. /lib/functions.sh

local enable
config_load amazon_ffs
config_get enable "ffs" "enable" "off"

[ "$ACTION" = "start" -a  "$enable" = "on" ] && {
	. /lib/ffs/ffs_core.sh
	ffs_rule_exit
	ffs_rule_init
}
