# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/accel-pptpd.list
/usr/sbin/pptpd
/lib/vpn/vpn_fc.sh
/usr/sbin/pptpctrl
/etc/init.d/pptpd
/lib/upgrade/keep.d/accel-pptpd
/etc/ppp/options.pptpd.sample
/etc/pptpd.conf
/usr/lib/pptpd/pptpd-logwtmp.so
/usr/sbin/bcrelay
/etc/config/pptpd
