# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/ffs.list
/etc/init.d/ffs
/lib/ffs/ffs_core.sh
/etc/hotplug.d/iface/66-ffs
/etc/init.d/ffs_load_config
/etc/hotplug.d/firewall/66-ffs
/lib/ffs/ffs_ratelimit.sh
/lib/ffs/ffs_tools.lua
/lib/ffs/ffs_network.sh
