# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/ledctrl.list
/etc/init.d/ledctrl
/usr/bin/ledcli
/etc/init.d/sys-led-end
/usr/bin/ledctrl
/etc/hotplug.d/usb/10-usb
/etc/config/ledctrl
/etc/init.d/sys-led-start
