# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/libcloud-api.list
/usr/bin/cloud-brd
/etc/hotplug.d/iface/98-cloud-brd
/etc/init.d/cloud_brd
/etc/cloud_config.cfg
/usr/lib/lua/cloud_dispatcher.lua
/usr/sbin/cloud_genSrvFile
/etc/certificate/2048_newroot.cer
