# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/smart-ip.list
/sbin/switch_mode_oldway
/usr/bin/smartipd
/etc/init.d/smartip
/etc/init.d/zzzapmode
/sbin/sync_mode_file
/sbin/restartlan
/sbin/remodel_ap_merge
/sbin/switch_mode
