# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/strongswan.list
/etc/init.d/ipsec
/usr/lib/ipsec/libstrongswan.so.0
/usr/lib/ipsec/libstrongswan.so.0.0.0
/etc/ipsec.secrets
/etc/strongswan.conf
/lib/upgrade/keep.d/strongswan
/usr/sbin/ipsec_client_tunnel_monitor
/etc/ipsec.conf.server
/etc/ipsec.conf.common
/etc/ipsec.user
