# Nicholas/TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/usr/lib/opkg/info/uhttpd.list
/usr/sbin/uhttpd
/etc/init.d/uhttpd
/etc/hotplug.d/iface/15-uhttpd
/etc/config/uhttpd
/lib/upgrade/keep.d/uhttpd
