# TPLINK_AX1800/_ax1800v1.20-up-ver1-3-2-P1[20210901-rel69355]_sign_2021-09-01_19.18.14.bin.extracted/ubifs-root/1972537161/rootfs_ubifs/etc/rc.d/S72usb_power
#!/bin/sh /etc/rc.common

START=72

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1

start() {
	insmod /lib/modules/*/extra/bcm_usb.ko
}

stop() {
	rmmod bcm_usb.ko
}

